# scripts/ai.cryptarch.trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.cryptarch.trader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/logannye/cryptarch/scripts/run_cryptarch.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/logannye/cryptarch</string>

    <!-- Auto-start on login + restart on crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restart cadence so we don't spin if the bot is broken. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Soft memory limit (8 GB). Hard cap via the OS would crash the bot;
         we just want a flag if it grows unexpectedly. -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>ResidentSetSize</key>
        <integer>8589934592</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/logannye/cryptarch/data/cryptarch_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/logannye/cryptarch/data/cryptarch_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
